# edp-monitor/rootfs/etc/services.d/edp-monitor/run
#!/usr/bin/with-contenv bashio
sleep 10
export DISPLAY=:99
bashio::log.info "Starting EDP Voucher Monitor..."
bashio::log.info "noVNC available at port 6080"
exec python3 /app/edp_monitor.py
